# sys-apps/busybox/files/ntpc.initd
#!/sbin/openrc-run
# -*- coding: utf-8; mode: sh; tab-width: 3; indent-tabs-mode: nil -*-
#
# Copyright 2018, 2020 Raffaello D. Di Napoli

# Before Gentoo bug #444718, there used to be a package with something similar to this.

name='NTP client'
command=busybox
procname=ntpd
command_args="${NTPD_OPTS}"
pidfile=/run/ntpd.pid

depend() {
   use dns logger net
}
